# buildlinux/debian/prerm
#!/bin/sh

if [ -x /bin/systemctl ]; then
    systemctl disable safejumperplus
    systemctl stop safejumperplus
fi

#DEBHELPER#
exit 0
