# buildlinux/debian/safejumperplus.upstart
# safejumper - Safejumper VPN Service.
#
# This task runs safejumper's service daemon.

description	"Safejumper Plus service daemon"

start on runlevel [2345]
stop on shutdown
expect daemon
respawn

script
    exec /opt/safejumperplus/safejumperplusservice
end script

